# jellyfin/etc/runit/sv/jellyfin/finish
#!/bin/sh
exec chpst -u soffiafdz podman stop -t 10 jellyfin
